# supervisord.conf
[supervisord]
nodaemon=true

[program:agent-api]
autorestart=true
startretries=10
command=python -m rl_engine.agent.api_server
directory=/app
stdout_logfile=/app/logs/agent_api.log
stderr_logfile=/app/logs/agent_api_error.log
autostart=true
autorestart=true