# supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# DQN Training Agent
[program:dqn-agent]
command=python -m rl_engine.agent.train_dqn
directory=/app
user=root
stdout_logfile=/app/logs/dqn_agent.log
stderr_logfile=/app/logs/dqn_agent_error.log
logfile_maxbytes=10MB
logfile_backups=5
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=60
environment=PYTHONPATH=/app,PYTHONUNBUFFERED=1,MLFLOW_TRACKING_URI=http://mlflow:5000,CI=false

# PPO Training Agent
[program:ppo-agent]
command=python -m rl_engine.agent.train_ppo
directory=/app
user=root
stdout_logfile=/app/logs/ppo_agent.log
stderr_logfile=/app/logs/ppo_agent_error.log
logfile_maxbytes=10MB
logfile_backups=5
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=60
environment=PYTHONPATH=/app,PYTHONUNBUFFERED=1,MLFLOW_TRACKING_URI=http://mlflow:5000,CI=false

# Process group
[group:rl-agents]
programs=dqn-agent,ppo-agent
priority=999